# deploy/memory-fs-backup.timer
[Unit]
Description=Daily memory-fs backup to GCS

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
